# rootfs/etc/supervisor.d/nginx.ini
[program:nginx]
command=/usr/local/bin/scythe_run_nginx
autostart=true
autorestart=true
stdout_logfile=/var/log/nginx.log
stderr_logfile=/var/log/nginx.error.log
logfile_maxbytes=1MB
logfile_backups=1
