# rootfs/etc/supervisor.d/sshd.ini
[program:sshd]
command=/usr/local/bin/scythe_run_sshd
autostart=true
autorestart=true
stdout_logfile=/var/log/sshd.log
stderr_logfile=/var/log/sshd.error.log
logfile_maxbytes=1MB
logfile_backups=1
