# rootfs/etc/supervisor.d/x11vnc.ini
[program:x11vnc]
command=/usr/bin/x11vnc -localhost -display :0
autostart=true
autorestart=true
stdout_logfile=/var/log/xvnc.log
stderr_logfile=/var/log/xvnc.error.log
logfile_maxbytes=1MB
logfile_backups=1
