# rootfs/etc/supervisor.d/xvfb.ini
[program:xvfb]
command=/usr/bin/xinit -- /usr/bin/Xvfb -screen 0 800x600x16 -ac
autostart=true
autorestart=true
stdout_logfile=/var/log/xvfb.log
stderr_logfile=/var/log/xvfb.error.log
logfile_maxbytes=1MB
logfile_backups=1
priority=1
